# agent/src/main/resources/rpm/nexis-agent.service
[Unit]
Description=Nexis Agent - System Metrics Collector
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/java -jar /opt/nexis-agent/Nexis_Agent.jar
WorkingDirectory=/opt/nexis-agent
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target